# dist/systemd/sstd.service
[Unit]
Description=BtrFS snapshot daemon.

[Service]
Type=simple
PIDFile=/tmp/sstd.pid
ExecStart=/usr/bin/sstd

[Install]
WantedBy=multi-user.target
