# src/makesnap@.timer
[Unit]
Description=Running MakeSnap on %I every 30min.

[Timer]
OnBootSec=1
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
